# root/etc/s6-overlay/s6-rc.d/p11-maven-customize/up
/etc/s6-overlay/s6-rc.d/p11-maven-customize/run
